# debian/postinst
#!/bin/sh -e

update-rc.d masysma-emergency-remote-control defaults
