# debian/postrm
#!/bin/sh -e

update-rc.d -f masysma-emergency-remote-control remove
